# contrib/init/bitcoind.conf
description "Bitcash Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BITCASHD_BIN="/usr/bin/bitcashd"
env BITCASHD_USER="bitcash"
env BITCASHD_GROUP="bitcash"
env BITCASHD_PIDDIR="/var/run/bitcashd"
# upstart can't handle variables constructed with other variables
env BITCASHD_PIDFILE="/var/run/bitcashd/bitcashd.pid"
env BITCASHD_CONFIGFILE="/etc/bitcash/bitcash.conf"
env BITCASHD_DATADIR="/var/lib/bitcashd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # bitcashd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BITCASHD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run bitcashd."
        echo "The setting must appear in $BITCASHD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcash Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BITCASHD_PIDDIR"
    chmod 0755 "$BITCASHD_PIDDIR"
    chown $BITCASHD_USER:$BITCASHD_GROUP "$BITCASHD_PIDDIR"
    chown $BITCASHD_USER:$BITCASHD_GROUP "$BITCASHD_CONFIGFILE"
    chmod 0660 "$BITCASHD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BITCASHD_PIDFILE" \
    --chuid $BITCASHD_USER:$BITCASHD_GROUP \
    --exec "$BITCASHD_BIN" \
    -- \
    -pid="$BITCASHD_PIDFILE" \
    -conf="$BITCASHD_CONFIGFILE" \
    -datadir="$BITCASHD_DATADIR" \
    -disablewallet \
    -daemon

